# mute-on-network.service
[Unit]
Description=Mute Audio on certain networks

[Service]
ExecStart=%h/Programme/scripts/netmute.py

[Install]
WantedBy=network.target
